# desktop_agent/docker-agent.service
[Unit]
Description=Docker Container Monitoring Agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=tiger
Group=docker
WorkingDirectory=/opt/docker-agent
Environment=PATH=/opt/docker-agent/venv/bin
ExecStart=/opt/docker-agent/venv/bin/python /opt/docker-agent/docker_agent.py
EnvironmentFile=/opt/docker-agent/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docker-agent

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/docker-agent
ReadWritePaths=/opt/docker-agent/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target